# bot-afiliados.service
[Unit]
Description=Bot Telegram Afiliados
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/bot-afiliados
ExecStart=/home/ubuntu/bot-afiliados/venv/bin/python bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Variáveis de ambiente carregadas do arquivo .env
EnvironmentFile=/home/ubuntu/bot-afiliados/.env

[Install]
WantedBy=multi-user.target
